# mcp-server/rootfs/etc/services.d/mcp-server/finish
#!/usr/bin/with-contenv bash
# ==============================================================================
# Take down the MCP Server service
# ==============================================================================

echo "[$(date +%H:%M:%S)] INFO: MCP Server service stopping..."

# Clean up any resources
if pgrep -f "node dist/index.js" > /dev/null 2>&1; then
    echo "[$(date +%H:%M:%S)] INFO: Stopping Node.js process..."
    pkill -TERM -f "node dist/index.js" || true
    sleep 2
    pkill -KILL -f "node dist/index.js" || true
fi

# Check exit code
if [[ "${1}" -eq 0 ]] || [[ "${1}" -eq 256 ]]; then
    echo "[$(date +%H:%M:%S)] INFO: MCP Server stopped cleanly"
else
    echo "[$(date +%H:%M:%S)] ERROR: MCP Server crashed with exit code ${1}" >&2
    # If it crashed, wait a bit before S6 restarts it
    sleep 5
fi