# mcp-server/rootfs/etc/services.d/mcp-server/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Start the MCP Server service
# ==============================================================================

# Get configuration from environment or use defaults
port="${PORT:-6789}"
log_level="${LOG_LEVEL:-info}"
connection_timeout="${CONNECTION_TIMEOUT:-30000}"
max_clients="${MAX_CLIENTS:-5}"
enable_entity_filtering="${ENABLE_ENTITY_FILTERING:-false}"

# Get HomeAssistant connection info
if [[ -z "${SUPERVISOR_TOKEN}" ]]; then
    echo "[$(date +%H:%M:%S)] ERROR: Supervisor token not found!" >&2
    echo "[$(date +%H:%M:%S)] INFO: Running in test mode"
    SUPERVISOR_TOKEN="test_token"
fi

export SUPERVISOR_TOKEN="${SUPERVISOR_TOKEN}"
export HASSIO_TOKEN="${SUPERVISOR_TOKEN}"
export HA_BASE_URL="http://supervisor/core"
export HOMEASSISTANT_URL="${HOMEASSISTANT_URL:-http://localhost:8123}"
export PORT="${port}"
export PORT_SSE="${port}"
export CONNECTION_MODE="sse"
export AUTHENTICATION_MODE="supervisor"
export LOG_LEVEL="${log_level}"
export CONNECTION_TIMEOUT="${connection_timeout}"
export MAX_CLIENTS="${max_clients}"
export ENABLE_ENTITY_FILTERING="${enable_entity_filtering}"

echo "[$(date +%H:%M:%S)] INFO: Starting MCP Server for Claude v1.2.0..."
echo "[$(date +%H:%M:%S)] INFO: Port: ${port}"
echo "[$(date +%H:%M:%S)] INFO: Log Level: ${log_level}"
echo "[$(date +%H:%M:%S)] INFO: Connection Mode: SSE (Server-Sent Events)"
echo "[$(date +%H:%M:%S)] INFO: Authentication: HomeAssistant Supervisor Token"
echo "[$(date +%H:%M:%S)] INFO: Max Clients: ${max_clients}"
echo "[$(date +%H:%M:%S)] INFO: Connection Timeout: ${connection_timeout}ms"

# Change to app directory
cd /app || exit 1

# Start the Node.js application
echo "[$(date +%H:%M:%S)] INFO: Starting Node.js application..."
exec node dist/index.js 2>&1